# data/supervisor/conf/crond.ini
[program:crond]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/crond -f
autostart=true
autorestart=true
user=root
redirect_stderr=true