# data/supervisor/conf/php-fpm.conf
[program:sshd]
command=/usr/local/sbin/php-fpm
